# stow/services/.config/systemd/user/backup.service
[Unit]
Description=Backup to GoogleDrive

[Service]
ExecStart=%h/.config/systemd/user/backup.sh

[Install]
WantedBy=default.target